# udev/95-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/chen/.Xauthority", RUN+="/usr/bin/xrandr --auto"
