# udev/95-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/chen/.Xauthority", RUN+="/usr/bin/timeout 3s /usr/bin/xrandr --auto"
