# udev/95-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/usr/bin/bash /home/chen/.open-display.sh"
